# scripts/watch_print_disable.service
[Unit]
Description=Watch Print Disable - Canon SELPHY CP1500 Auto-heal
After=multi-user.target cups.service
Requires=cups.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/watch_print_disable.sh
User=root
StandardOutput=journal
StandardError=journal
StandardInput=null

[Install]
WantedBy=multi-user.target